# Truehire/Truehire/deploy/systemd/truehire-backend.service
[Unit]
Description=TrueHire backend
After=network.target

[Service]
Type=simple
User=ubuntu
WorkingDirectory=/opt/truehire/Truehire/backend
EnvironmentFile=/etc/truehire/backend.env
ExecStart=/usr/bin/npm start
Restart=always
RestartSec=5

[Install]
WantedBy=multi-user.target

